# templates/default/kibana.init.erb
#!/bin/sh -e
### BEGIN INIT INFO
# Provides:          kibana
# Required-Start:    $remote_fs $syslog
# Required-Stop:     $remote_fs $syslog
# Default-Start:     2 3 4 5
# Default-Stop:      0 1 6
# Short-Description: Make sense of a mountain of logs.
### END INIT INFO

#
# chkconfig: 2345 20 80
#
KIBANA_HOME=<%= @kibana_home %>
USER=<%= @user %>

KIBANA_CMD="cd $KIBANA_HOME && bundle exec ruby kibana-daemon.rb $1"

case "$1" in
start)
  su - $USER -c "$KIBANA_CMD"
  RETVAL=$?
;;

stop)
  su - $USER -c "$KIBANA_CMD"
  RETVAL=$?
;;

restart)
  su - $USER -c "$KIBANA_CMD"
  RETVAL=$?
;;

status)
  su - $USER -c "$KIBANA_CMD"
;;

force-reload)
  su - $USER -c "cd $KIBANA_HOME && bundle exec kibana-daemon.rb zap"
  su - $USER -c "cd $KIBANA_HOME && bundle exec kibana-daemon.rb start"
  RETVAL=$?
;;

*)
  echo "Usage: $0 {start|stop|restart|status|force-reload}"
  exit 1
esac

exit $RETVAL
