# templates/default/logstash_beaver.conf.erb
description "Logstash beaver"
author "Chef"

start on (filesystem and net-device-up)
stop on runlevel [!2345]

respawn
respawn limit 5 30

<% if @supports_setuid  && !@ensure_supplemental_groups %>
setuid <%= @user %>
<% unless @group.to_s.empty? -%>
setgid <%= @group %>
<% end -%>
<% end %>
chdir <%= node['logstash']['basedir'] %>/beaver

console output

<% if @supports_setuid && @ensure_supplemental_groups %>
exec sudo -u <%= @user %> <%= @cmd %> >> <%= @log %> 2>&1
<% else %>
exec <%= @cmd %> >> <%= @log %> 2>&1
<% end %>

emits logstash-beaver-running
