# templates/default/logstash_agent.conf.erb
description "Logstash agent"
author "Chef"

start on (filesystem and net-device-up)
stop on runlevel [!2345]

respawn
respawn limit 5 30

chdir <%= node['logstash']['basedir'] %>/agent
setuid <%= node['logstash']['user'] %>

script
  export LOGSTASH_HOME="<%= node['logstash']['basedir'] %>/agent"
  export HOME=$LOGSTASH_HOME
  export GC_OPTS="<%= node['logstash']['agent']['gc_opts'] %>"
  export JAVA_OPTS="-server -Xms<%= node['logstash']['agent']['xms'] %> -Xmx<%= node['logstash']['agent']['xmx'] %> -Djava.io.tmpdir=$LOGSTASH_HOME/tmp/ <%= node[:logstash][:agent][:java_opts] %> <%= '-Djava.net.preferIPv4Stack=true' if node[:logstash][:agent][:ipv4_only] %>"
  export LOGSTASH_OPTS="agent -f $LOGSTASH_HOME/etc/shipper.conf -l $LOGSTASH_HOME/log/logstash.log"
  <% if node['logstash']['agent']['debug'] -%>
  export LOGSTASH_OPTS="$LOGSTASH_OPTS -vv"
  <% end -%>
  export OPTS="$JAVA_OPTS $GC_OPTS -jar $LOGSTASH_HOME/lib/logstash.jar $LOGSTASH_OPTS"

  exec /usr/bin/java $OPTS
end script

emits logstash-agent-running

