# templates/default/logstash_web.conf.erb
description "Logstash"
author "Chef"

start on (filesystem and net-device-up)
stop on runlevel [!2345]

respawn
respawn limit 5 30
limit nofile 65550 65550

chdir <%= node['logstash']['server']['home'] %>
setuid <%= node['logstash']['user'] %>

script
  export LOGSTASH_HOME="<%= node['logstash']['server']['home'] %>"
  export HOME=$LOGSTASH_HOME
  export GC_OPTS="<%= node['logstash']['server']['gc_opts'] %>"
  export JAVA_OPTS="-server -Xms<%= node['logstash']['server']['xms'] %> -Xmx<%= node['logstash']['server']['xmx'] %> -Djava.io.tmpdir=$LOGSTASH_HOME/tmp/ <%= node['logstash']['server']['java_opts'] %> <%= '-Djava.net.preferIPv4Stack=true' if node['logstash']['agent']['ipv4_only'] %>"
  export LOGSTASH_OPTS="web -a <%= node['logstash']['server']['web']['address'] %> -p <%= node['logstash']['server']['web']['port'] %>"
  export OPTS="$JAVA_OPTS $GC_OPTS -jar $LOGSTASH_HOME/lib/logstash.jar $LOGSTASH_OPTS"

  exec /usr/bin/java $OPTS
end script

emits logstash-web-running
