# templates/default/logstash_server.conf.erb
description "Logstash"
author "Chef"

start on (filesystem and net-device-up)
stop on runlevel [!2345]

respawn
respawn limit 5 30
limit nofile 65550 65550

chdir <%= node['logstash']['server']['home'] %>

<% unless node['logstash']['server']['upstart_with_sudo'] -%>
setuid <%= node['logstash']['user'] %>
  <% unless node['logstash']['supervisor_gid'].to_s.empty? -%>
setgid <%= node['logstash']['supervisor_gid'] %>
  <% end -%>
<% end -%>

script
  export LOGSTASH_HOME="<%= node['logstash']['server']['home'] %>"
  export HOME=$LOGSTASH_HOME
  export GC_OPTS="<%= node['logstash']['server']['gc_opts'] %>"
  export JAVA_OPTS="-server -Xms<%= node['logstash']['server']['xms'] %> -Xmx<%= node['logstash']['server']['xmx'] %> -Djava.io.tmpdir=$LOGSTASH_HOME/tmp/ <%= node['logstash']['server']['java_opts'] %> <%= '-Djava.net.preferIPv4Stack=true' if node['logstash']['agent']['ipv4_only'] %>"
  export LOGSTASH_OPTS="agent -f <%= node['logstash']['server']['home'] %>/etc/conf.d"
  <% if node['logstash']['server']['pluginpath'] -%>
  export LOGSTASH_OPTS="$LOGSTASH_OPTS --pluginpath <%= node['logstash']['server']['pluginpath'] %>"
  <% end -%>
  <% if node['logstash']['server']['workers'] %>
  export LOGSTASH_OPTS="$LOGSTASH_OPTS -w <%= node['logstash']['server']['workers'] %>"
  <% end %>
  <% if node['logstash']['server']['log_file'] -%>
  export LOGSTASH_OPTS="$LOGSTASH_OPTS -l <%= node['logstash']['server']['home'] %>/log/<%= node['logstash']['server']['log_file'] %>"
  <% end -%>
  <% if node['logstash']['server']['debug'] -%>
  export LOGSTASH_OPTS="$LOGSTASH_OPTS -vv"
  <% end -%>
  export OPTS="$JAVA_OPTS $GC_OPTS -jar $LOGSTASH_HOME/lib/logstash.jar $LOGSTASH_OPTS"

  <% if node['logstash']['server']['upstart_with_sudo'] -%>
  exec sudo -u <%= node['logstash']['user'] %> /usr/bin/java $OPTS
  <% else -%>
  exec /usr/bin/java $OPTS
  <% end -%>
end script

emits logstash-server-running
