# templates/default/logstash_agent.conf.erb
description "Logstash agent"
author "Chef"

start on (filesystem and net-device-up)
stop on runlevel [!2345]
limit nofile 65550 65550

respawn
respawn limit <%=node['logstash']['agent']['upstart_respawn_count']%> <%=node['logstash']['agent']['upstart_respawn_timeout']%>

chdir <%= node['logstash']['agent']['home'] %>

<% unless node['logstash']['agent']['upstart_with_sudo'] -%>
setuid <%= node['logstash']['user'] %>
  <% unless node['logstash']['supervisor_gid'].to_s.empty? -%>
setgid <%= node['logstash']['supervisor_gid'] %>
  <% end -%>
<% end -%>

script
  export LOGSTASH_HOME="<%= node['logstash']['agent']['home'] %>"
  export HOME=$LOGSTASH_HOME
  export GC_OPTS="<%= node['logstash']['agent']['gc_opts'] %>"
  export JAVA_OPTS="-server -Xms<%= node['logstash']['agent']['xms'] %> -Xmx<%= node['logstash']['agent']['xmx'] %> -Djava.io.tmpdir=$LOGSTASH_HOME/tmp/ <%= node['logstash']['agent']['java_opts'] %> <%= '-Djava.net.preferIPv4Stack=true' if node['logstash']['agent']['ipv4_only'] %>"
  export LOGSTASH_OPTS="agent -f <%= node['logstash']['agent']['config_dir'] %>"
  <% if node['logstash']['agent']['workers'] %>
  export LOGSTASH_OPTS="$LOGSTASH_OPTS -w <%= node['logstash']['agent']['workers'] %>"
  <% end %>
  <% if node['logstash']['agent']['debug'] -%>
  export LOGSTASH_OPTS="$LOGSTASH_OPTS -vv"
  <% end -%>
  <% if node['logstash']['agent']['log_file'] -%>
  export LOGSTASH_OPTS="$LOGSTASH_OPTS -l <%= node['logstash']['agent']['log_file'] %>"
  <% end -%>

  export OPTS="$JAVA_OPTS $GC_OPTS -jar $LOGSTASH_HOME/lib/logstash.jar $LOGSTASH_OPTS"

  <% if node['logstash']['agent']['upstart_with_sudo'] -%>
  exec sudo -u <%= node['logstash']['user'] %> /usr/bin/java $OPTS
  <% else -%>
  exec /usr/bin/java $OPTS
  <% end -%>
end script

emits logstash-agent-running
