# templates/default/logstash_beaver.conf.erb
description "Logstash beaver"
author "Chef"

start on (filesystem and net-device-up)
stop on runlevel [!2345]

respawn
respawn limit 5 30

<% if @supports_setuid %>
setuid <%= @user %>
<% unless @group.to_s.empty? -%>
setgid <%= @group %>
<% end -%>
<% end %>
chdir <%= node['logstash']['basedir'] %>/beaver

console output

exec <%= @cmd %> >> <%= @log %> 2>&1

emits logstash-beaver-running
